# src/crontab
@reboot    /renew.sh
0 3 * * *  /renew.sh